# lesson_2/etc/cron.d
# Every day 5 a.m.
0 5 * * * /home/ubuntu/scripts/lesson_2/bin/lesson_2
